# myservice.timer
[Unit]
Description=My timer
Requires=myservice.service

[Timer]
Unit=myMonitor.service
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
